# monitoring/systemd/redisforge-metrics-push.service
[Unit]
Description=RedisForge Metrics Push Service
After=docker.service
Requires=docker.service

[Service]
Type=simple
# For Amazon Linux 2023, use: User=ec2-user
# For Ubuntu 24.04 LTS, use: User=ubuntu
User=ec2-user
# For Amazon Linux 2023, use: WorkingDirectory=/home/ec2-user/RedisForge
# For Ubuntu 24.04 LTS, use: WorkingDirectory=/home/ubuntu/RedisForge
WorkingDirectory=/home/ec2-user/RedisForge
ExecStart=/home/ec2-user/RedisForge/scripts/push-metrics.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=redisforge-metrics

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=multi-user.target
